# deploy/systemd/arabesque-report-weekly.timer
[Unit]
Description=Arabesque — rapport hebdomadaire (dimanche 20h UTC)

[Timer]
OnCalendar=Sun *-*-* 20:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
