# deploy/systemd/installed/arabesque-telegram-bot.service
[Unit]
Description=Arabesque — bot Telegram interactif (lecture)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/var/home/machine/dev/arabesque
ExecStart=/var/home/machine/dev/arabesque/.venv/bin/python -m arabesque.bot.telegram_bot
Restart=on-failure
RestartSec=15
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
